# griddle-monitor.service
[Unit]
Description=Griddle Temperature Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/griddle-monitor
ExecStart=/home/pi/griddle-monitor/venv/bin/python /home/pi/griddle-monitor/app.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
